# fedora-solr-gsearch/services/tomcat/run
#!/usr/bin/env bash

${FEDORA_HOME}/tomcat/bin/catalina.sh run
